# scripts/carousel/crontab
# Fetch espacios disponibles every hour during working hours
0 6-15 * * 1-5 /app/scripts/carousel/fetch_espacios.sh
